# deploy/systemd/xresidual-logger-orderbooks.timer
[Unit]
Description=xResidual order-book logger every 30 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
