# deploy/systemd/xresidual-logger-oddsapi.service
[Unit]
Description=xResidual bookmaker odds logger (The Odds API)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=azureuser
WorkingDirectory=/home/azureuser/xResidual/logger
ExecStart=/home/azureuser/xResidual/.venv/bin/python run.py --venues oddsapi
StandardOutput=append:/home/azureuser/xResidual/logger/data/svc.log
StandardError=append:/home/azureuser/xResidual/logger/data/svc.log
