# deploy/systemd/xresidual-tape-cleanup.timer
[Unit]
Description=Run the xResidual tape cleanup daily (keep disk in check over the tournament)

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
