# deploy/systemd/xresidual-site.timer
[Unit]
Description=xResidual live-site refresh every 30 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
